# deploy/aps-etl.timer
[Unit]
Description=Executa ETL APS todos os dias às 04:00

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
